# docker/php/laravel.cron
* * * * * cd /var/www && php artisan subastas:cerrar >> /var/log/cron.log 2>&1
